# ip-tracker.cron
0 * * * * root /opt/ip-tracker/ip-tracker.sh remote.example.com
