# ip-tracker.cron
0 * * * * root /data/ip-tracker/tracker.sh lab.jasapple.com
